# profile/airootfs/etc/skel/.config/autostart/ksuperkey-right.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=Ksuperkey Right
Comment=Right Super Key
Exec=ksuperkey -e 'Super_R=Alt_L|F1'
